# etc/supervisor/torrentpotato.t411.conf
[program:torrentpotato.t411-gunicorn-prod]
environment=PROD=True
command=/var/www/torrentpotato.t411/prod/venv/bin/gunicorn -c /var/www/torrentpotato.t411/prod/current/etc/gunicorn.conf.py start:app
user=fabric
autostart=true
autorestart=true
redirect_stderr=true
killasgroup=true
stopasgroup=true
